# dyndns-update.service
[Unit]
Description=update dydns via docker

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run --rm dyndns-update
